# docker/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm
user=www
startsecs=0
numprocs=1
autostart=true
autorestart=true
priority=100
process_name=%(program_name)s_%(process_num)02d
