# docker/supervisor/conf.d/worker.conf
[program:worker]
command=php /var/www/app/worker.php
user=www
numprocs=50
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d